# docker_build/oracle19/0.0.4/api.ini
[program:api]
directory=/api/web
command=/bin/startApi.sh
environment=LC_ALL=en_US.utf8,LANG=en_US.utf8
user=root
autorestart=false
stdout_logfile=/tmp/api.log
stderr_logfile=/tmp/api.log
stopsignal=TERM
