# docker_build/oracle19/0.0.5/oracle.ini
[program:oracle]
command=/bin/startOracle.sh 
user=oracle
autorestart=false
stdout_logfile=/tmp/startoracle.log
stderr_logfile=/tmp/startoracle.log
stopsignal=TERM
